# dist/supervisord.conf
[supervisord]
nodaemon=true

[program:tarantool]
command=/usr/local/bin/tarantool /opt/tarantool/docker.lua

[program:nginx]
command=nginx -g "daemon off;"
